# .config/systemd/user/backup.timer
[Unit]
Description=Backup timer

[Timer]
OnBootSec=hourly
OnUnitInactiveSec=hourly

[Install]
WantedBy=default.target
