# .config/systemd/user/backup.service
[Unit]
Description=System backup


[Service]
CPUQuota=10%
Type=simple
ExecStart=%h/.local/bin/backup home

[Install]
WantedBy=default.target
