# .config/systemd/user/backup.service
[Unit]
Description=Borg backup


[Service]
CPUQuota=60%
Type=simple
ExecStart=%h/.local/bin/backup home

[Install]
WantedBy=default.target
