# .config/systemd/user/mail.timer
[Unit]
Description=Mail target timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=1m
Unit=mail.target

[Install]
WantedBy=default.target
