# .config/systemd/user/mail.target
[Unit]
Description=Mail target
StopWhenUnneeded=true
