# .config/systemd/user/ssh-tpm-agent.service.d/override.conf
#
[Service]
ExecStart=
ExecStart=/usr/bin/ssh-tpm-agent -A /run/user/1000/ssh-agent.socket
